# letsencrypt/letsencrypt-renewal.timer
[Unit]
Description=Weekly SSL cert check/renewal

[Timer]
OnCalendar=Thu 00:03
Persistent=true

[Install]
WantedBy=timers.target
